# scripts/backend.service
[Unit]
Description=TravelForge API Server
After=network.target

[Service]
Type=simple
User=user
WorkingDirectory=/usr/share/nginx/html/apps/TravelFrog/backend
Environment=NODE_ENV=production
EnvironmentFile=/usr/share/nginx/html/apps/TravelFrog/backend/.env
ExecStart=/usr/bin/node /usr/share/nginx/html/apps/TravelFrog/backend/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

